# files/etc/systemd/system/systemd-journald.service.d/override.conf
#reconfigure systemd-journald.service to add /var hook problem solver 
[Unit]
RequiresMountsFor=/var/log/journal
